# var/lib/dpkg/info/samba.conffiles
/etc/cron.daily/samba
/etc/init.d/nmbd
/etc/init.d/samba-ad-dc
/etc/init.d/smbd
/etc/logrotate.d/samba
/etc/ufw/applications.d/samba
