# var/lib/dpkg/info/samba.list
/.
/etc
/etc/apparmor.d
/etc/apparmor.d/samba
/etc/cron.daily
/etc/cron.daily/samba
/etc/init.d
/etc/init.d/nmbd
/etc/init.d/samba-ad-dc
/etc/init.d/smbd
/etc/logrotate.d
/etc/logrotate.d/samba
/etc/ufw
/etc/ufw/applications.d
/etc/ufw/applications.d/samba
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/nmbd.service
/lib/systemd/system/samba-ad-dc.service
/lib/systemd/system/smbd.service
/usr
/usr/bin
/usr/bin/dumpmscat
/usr/bin/mvxattr
/usr/bin/oLschema2ldif
/usr/bin/pdbedit
/usr/bin/profiles
/usr/bin/sharesec
/usr/bin/smbcontrol
/usr/bin/smbstatus
/usr/lib
/usr/lib/arm-linux-gnueabihf
/usr/lib/arm-linux-gnueabihf/samba
/usr/lib/arm-linux-gnueabihf/samba/libHDB-SAMBA4.so.0
/usr/lib/arm-linux-gnueabihf/samba/libdb-glue.so.0
/usr/lib/arm-linux-gnueabihf/samba/libhdb-samba4.so.11.0.2
/usr/lib/arm-linux-gnueabihf/samba/libkdc-samba4.so.2.0.0
/usr/lib/arm-linux-gnueabihf/samba/libmscat.so.0
/usr/lib/arm-linux-gnueabihf/samba/libpac.so.0
/usr/lib/arm-linux-gnueabihf/samba/service
/usr/lib/arm-linux-gnueabihf/samba/service/cldap.so
/usr/lib/arm-linux-gnueabihf/samba/service/dcerpc.so
/usr/lib/arm-linux-gnueabihf/samba/service/dns.so
/usr/lib/arm-linux-gnueabihf/samba/service/dns_update.so
/usr/lib/arm-linux-gnueabihf/samba/service/drepl.so
/usr/lib/arm-linux-gnueabihf/samba/service/kcc.so
/usr/lib/arm-linux-gnueabihf/samba/service/kdc.so
/usr/lib/arm-linux-gnueabihf/samba/service/ldap.so
/usr/lib/arm-linux-gnueabihf/samba/service/nbtd.so
/usr/lib/arm-linux-gnueabihf/samba/service/ntp_signd.so
/usr/lib/arm-linux-gnueabihf/samba/service/s3fs.so
/usr/lib/arm-linux-gnueabihf/samba/service/winbindd.so
/usr/lib/arm-linux-gnueabihf/samba/service/wrepl.so
/usr/lib/python3
/usr/lib/python3/dist-packages
/usr/lib/python3/dist-packages/samba
/usr/lib/python3/dist-packages/samba/dckeytab.cpython-38-arm-linux-gnueabihf.so
/usr/sbin
/usr/sbin/eventlogadm
/usr/sbin/mksmbpasswd
/usr/sbin/nmbd
/usr/sbin/samba
/usr/sbin/samba-gpupdate
/usr/sbin/samba_dnsupdate
/usr/sbin/samba_downgrade_db
/usr/sbin/samba_spnupdate
/usr/sbin/samba_upgradedns
/usr/sbin/smbd
/usr/share
/usr/share/bug
/usr/share/bug/samba
/usr/share/doc
/usr/share/doc/samba
/usr/share/doc/samba/README.Debian
/usr/share/doc/samba/TODO.Debian
/usr/share/doc/samba/copyright
/usr/share/doc/samba/examples
/usr/share/doc/samba/examples/LDAP
/usr/share/doc/samba/examples/LDAP/README
/usr/share/doc/samba/examples/LDAP/get_next_oid
/usr/share/doc/samba/examples/LDAP/ol-schema-migrate.pl
/usr/share/doc/samba/examples/LDAP/samba-nds.schema
/usr/share/doc/samba/examples/LDAP/samba-schema-FDS.ldif
/usr/share/doc/samba/examples/LDAP/samba-schema-netscapeds5.x.README
/usr/share/doc/samba/examples/LDAP/samba-schema.IBMSecureWay
/usr/share/doc/samba/examples/LDAP/samba.ldif
/usr/share/doc/samba/examples/LDAP/samba.schema
/usr/share/doc/samba/examples/LDAP/samba.schema.at.IBM-DS
/usr/share/doc/samba/examples/LDAP/samba.schema.oc.IBM-DS
/usr/share/doc/samba/examples/logon
/usr/share/doc/samba/examples/logon/genlogon
/usr/share/doc/samba/examples/logon/genlogon/genlogon.pl
/usr/share/doc/samba/examples/logon/mklogon
/usr/share/doc/samba/examples/logon/mklogon/mklogon.conf
/usr/share/doc/samba/examples/logon/mklogon/mklogon.pl
/usr/share/doc/samba/examples/logon/ntlogon
/usr/share/doc/samba/examples/logon/ntlogon/README
/usr/share/doc/samba/examples/logon/ntlogon/ntlogon.conf
/usr/share/doc/samba/examples/logon/ntlogon/ntlogon.py
/usr/share/doc/samba/examples/printing
/usr/share/doc/samba/examples/printing/VampireDriversFunctions
/usr/share/doc/samba/examples/printing/prtpub.c
/usr/share/doc/samba/examples/printing/readme.prtpub
/usr/share/doc/samba/examples/printing/smbprint.sysv
/usr/share/doc/samba/examples/smbadduser.in
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/samba
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/log2pcap.1.gz
/usr/share/man/man1/mvxattr.1.gz
/usr/share/man/man1/oLschema2ldif.1.gz
/usr/share/man/man1/profiles.1.gz
/usr/share/man/man1/sharesec.1.gz
/usr/share/man/man1/smbcontrol.1.gz
/usr/share/man/man1/smbstatus.1.gz
/usr/share/man/man8
/usr/share/man/man8/eventlogadm.8.gz
/usr/share/man/man8/idmap_script.8.gz
/usr/share/man/man8/mksmbpasswd.8.gz
/usr/share/man/man8/nmbd.8.gz
/usr/share/man/man8/pdbedit.8.gz
/usr/share/man/man8/samba-gpupdate.8.gz
/usr/share/man/man8/samba.8.gz
/usr/share/man/man8/samba_downgrade_db.8.gz
/usr/share/man/man8/smbd.8.gz
/usr/share/samba
/usr/share/samba/admx
/usr/share/samba/admx/en-US
/usr/share/samba/admx/en-US/samba.adml
/usr/share/samba/admx/samba.admx
/usr/share/samba/mdssvc
/usr/share/samba/mdssvc/elasticsearch_mappings.json
/usr/share/samba/setup
/usr/share/samba/setup/ad-schema
/usr/share/samba/setup/ad-schema/AD_DS_Attributes__Windows_Server_2012_R2.ldf
/usr/share/samba/setup/ad-schema/AD_DS_Attributes__Windows_Server_2016.ldf
/usr/share/samba/setup/ad-schema/AD_DS_Classes__Windows_Server_2012_R2.ldf
/usr/share/samba/setup/ad-schema/AD_DS_Classes__Windows_Server_2016.ldf
/usr/share/samba/setup/ad-schema/Attributes_for_AD_DS__Windows_Server_2008_R2.ldf
/usr/share/samba/setup/ad-schema/Attributes_for_AD_DS__Windows_Server_2012.ldf
/usr/share/samba/setup/ad-schema/Classes_for_AD_DS__Windows_Server_2008_R2.ldf
/usr/share/samba/setup/ad-schema/Classes_for_AD_DS__Windows_Server_2012.ldf
/usr/share/samba/setup/ad-schema/MS-AD_Schema_2K8_Attributes.txt
/usr/share/samba/setup/ad-schema/MS-AD_Schema_2K8_Classes.txt
/usr/share/samba/setup/ad-schema/MS-AD_Schema_2K8_R2_Attributes.txt
/usr/share/samba/setup/ad-schema/MS-AD_Schema_2K8_R2_Classes.txt
/usr/share/samba/setup/adprep
/usr/share/samba/setup/adprep/WindowsServerDocs
/usr/share/samba/setup/adprep/WindowsServerDocs/Forest-Wide-Updates.md
/usr/share/samba/setup/adprep/WindowsServerDocs/Sch49.ldf.diff
/usr/share/samba/setup/adprep/WindowsServerDocs/Sch50.ldf.diff
/usr/share/samba/setup/adprep/WindowsServerDocs/Sch51.ldf.diff
/usr/share/samba/setup/adprep/WindowsServerDocs/Sch57.ldf.diff
/usr/share/samba/setup/adprep/WindowsServerDocs/Sch59.ldf.diff
/usr/share/samba/setup/adprep/WindowsServerDocs/Schema-Updates.md
/usr/share/samba/setup/adprep/fix-forest-rev.ldf
/usr/share/samba/setup/aggregate_schema.ldif
/usr/share/samba/setup/display-specifiers
/usr/share/samba/setup/display-specifiers/DisplaySpecifiers-Win2k0.txt
/usr/share/samba/setup/display-specifiers/DisplaySpecifiers-Win2k3.txt
/usr/share/samba/setup/display-specifiers/DisplaySpecifiers-Win2k3R2.txt
/usr/share/samba/setup/display-specifiers/DisplaySpecifiers-Win2k8.txt
/usr/share/samba/setup/display-specifiers/DisplaySpecifiers-Win2k8R2.txt
/usr/share/samba/setup/dns_update_list
/usr/share/samba/setup/extended-rights.ldif
/usr/share/samba/setup/idmap_init.ldif
/usr/share/samba/setup/krb5.conf
/usr/share/samba/setup/named.conf
/usr/share/samba/setup/named.conf.dlz
/usr/share/samba/setup/named.conf.update
/usr/share/samba/setup/named.txt
/usr/share/samba/setup/prefixMap.txt
/usr/share/samba/setup/provision.ldif
/usr/share/samba/setup/provision.reg
/usr/share/samba/setup/provision.zone
/usr/share/samba/setup/provision_basedn.ldif
/usr/share/samba/setup/provision_basedn_modify.ldif
/usr/share/samba/setup/provision_basedn_options.ldif
/usr/share/samba/setup/provision_basedn_references.ldif
/usr/share/samba/setup/provision_computers_add.ldif
/usr/share/samba/setup/provision_computers_modify.ldif
/usr/share/samba/setup/provision_configuration.ldif
/usr/share/samba/setup/provision_configuration_basedn.ldif
/usr/share/samba/setup/provision_configuration_modify.ldif
/usr/share/samba/setup/provision_configuration_references.ldif
/usr/share/samba/setup/provision_dns_accounts_add.ldif
/usr/share/samba/setup/provision_dns_add_samba.ldif
/usr/share/samba/setup/provision_dnszones_add.ldif
/usr/share/samba/setup/provision_dnszones_modify.ldif
/usr/share/samba/setup/provision_dnszones_partitions.ldif
/usr/share/samba/setup/provision_group_policy.ldif
/usr/share/samba/setup/provision_init.ldif
/usr/share/samba/setup/provision_partitions.ldif
/usr/share/samba/setup/provision_privilege.ldif
/usr/share/samba/setup/provision_rootdse_add.ldif
/usr/share/samba/setup/provision_rootdse_modify.ldif
/usr/share/samba/setup/provision_schema_basedn.ldif
/usr/share/samba/setup/provision_schema_basedn_modify.ldif
/usr/share/samba/setup/provision_self_join.ldif
/usr/share/samba/setup/provision_self_join_config.ldif
/usr/share/samba/setup/provision_self_join_modify.ldif
/usr/share/samba/setup/provision_self_join_modify_config.ldif
/usr/share/samba/setup/provision_self_join_modify_schema.ldif
/usr/share/samba/setup/provision_users.ldif
/usr/share/samba/setup/provision_users_add.ldif
/usr/share/samba/setup/provision_users_modify.ldif
/usr/share/samba/setup/provision_well_known_sec_princ.ldif
/usr/share/samba/setup/schema_samba4.ldif
/usr/share/samba/setup/secrets.ldif
/usr/share/samba/setup/secrets_dns.ldif
/usr/share/samba/setup/secrets_init.ldif
/usr/share/samba/setup/share.ldif
/usr/share/samba/setup/spn_update_list
/usr/share/samba/setup/ypServ30.ldif
/usr/share/samba/update-apparmor-samba-profile
/var
/var/lib
/var/lib/samba
/var/lib/samba/printers
/var/lib/samba/printers/COLOR
/var/lib/samba/printers/IA64
/var/lib/samba/printers/W32ALPHA
/var/lib/samba/printers/W32MIPS
/var/lib/samba/printers/W32PPC
/var/lib/samba/printers/W32X86
/var/lib/samba/printers/WIN40
/var/lib/samba/printers/x64
/var/spool
/var/spool/samba
/usr/lib/arm-linux-gnueabihf/samba/libhdb-samba4.so.11
/usr/lib/arm-linux-gnueabihf/samba/libkdc-samba4.so.2
/usr/share/bug/samba/presubj
/usr/share/bug/samba/script
/usr/share/doc/samba/NEWS.Debian.gz
/usr/share/doc/samba/changelog.Debian.gz
